# SRCD/autorun/99_detect_new_disk.rules
SUBSYSTEM=="block", ACTION=="add", RUN+="/root/add_disk_to_wipe.sh %N"
